# launchd/com.fernandozamora.polymarket-insider-finder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.fernandozamora.polymarket-insider-finder</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/fernandozamora/.pyenv/versions/3.9.2/bin/python3</string>
		<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/polymarket_insider_finder.py</string>
		<string>--service</string>
		<string>--telegram</string>
		<string>--telegram-env-file</string>
		<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/config/telegram.env</string>
		<string>--interval</string>
		<string>60</string>
		<string>--top</string>
		<string>10</string>
		<string>--limit-per-page</string>
		<string>100</string>
		<string>--db-path</string>
		<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/data/polymarket_insider.sqlite3</string>
		<string>--rules-file</string>
		<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/config/signal_rules.json</string>
		<string>--log-file</string>
		<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/logs/polymarket_insider_finder.log</string>
		<string>--notify-top</string>
		<string>3</string>
		<string>--notification-cooldown</string>
		<string>1800</string>
		<string>--min-liquidity</string>
		<string>2000.0</string>
		<string>--min-volume-24h</string>
		<string>250.0</string>
		<string>--min-oi-abs</string>
		<string>5000.0</string>
		<string>--min-oi-pct</string>
		<string>0.04</string>
		<string>--min-price-move</string>
		<string>0.06</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/fernandozamora/Desktop/polymarket_insider_finder</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/logs/launchd.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/fernandozamora/Desktop/polymarket_insider_finder/logs/launchd.stderr.log</string>
</dict>
</plist>
